# deploy/systemd/spotinotifs.service
[Unit]
Description=SpotiNotifs web service (Docker Compose)
Requires=docker.service
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=@PROJECT_DIR@
ExecStart=/usr/bin/docker compose up --no-build server
ExecStop=/usr/bin/docker compose down
TimeoutStartSec=0
TimeoutStopSec=90
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
